# deploy/systemd/diet-restore-test.timer
[Unit]
Description=Run diet-restore-test.service nightly at 03:45

[Timer]
OnCalendar=*-*-* 03:45:00
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
