# deploy/systemd/diet-app.service
[Unit]
Description=Diet App (Gunicorn + Uvicorn workers)
After=network.target

[Service]
WorkingDirectory=/home/dx/diet-app
Environment=PATH=/home/dx/diet-app/.venv/bin:/usr/bin
ExecStart=/home/dx/diet-app/.venv/bin/gunicorn app.main:app --workers 2 --worker-class uvicorn.workers.UvicornWorker --bind 127.0.0.1:8010 --pid /home/dx/diet-app/gunicorn.pid --graceful-timeout 30 --timeout 60 --log-level info
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=5
KillMode=mixed
TimeoutStopSec=60

[Install]
WantedBy=default.target
