# deploy/systemd/overrides/diet-healthcheck.service.d/override.conf
[Service]
Environment=AUTO_RESTART=yes
Environment=RESTART_AFTER_FAILS=3
